# docs/Archived/crontab/.crontab.conf
# Date: 2021-07-14 00:17:22
# File: crontab.conf
0 21 */1 * * bash ~/dotfiles/archive/sh/git-auto.sh
